# Backend/deploy/systemd/araneae-control.service
[Unit]
Description=Araneae GoRefactor Control Node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=araneae
Group=araneae
WorkingDirectory=/opt/araneae/GoRefactor
EnvironmentFile=/etc/araneae/control.env
ExecStart=/opt/araneae/GoRefactor/bin/araneae-control
Restart=always
RestartSec=2
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
